# data/kurs/crontab
* 12 * * * pobierz_A.php
* 12 * * 3 pobierz_B.php
* 8 * * * pobierz_C.php
