# data/lotto/crontab
* 23 * * * pobierz.php
* 0 * * * pobierz.php
* 15 * * * pobierz.php
* 16 * * * pobierz.php
